# supervisord.conf
[supervisord]
nodaemon=true

[program:httpd]
command=/etc/init.d/httpd start

[program:php-fpm]
command=/etc/init.d/php-fpm start
